# atlas-es-indexer.target
[Unit]
Description=Atlas measurement indexer Target
